# services/user-service/nexabank-user-service.service
[Unit]
Description=NexaBank User Service
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=120
StartLimitBurst=3

[Service]
Type=simple
User=nexabank
Group=nexabank
WorkingDirectory=/opt/nexabank/user-service

# The + prefix means: run this as root (needed to write /etc/nexabank/*.env)
# even though the main process runs as nexabank
ExecStartPre=+/bin/bash /opt/nexabank/user-service/bootstrap.sh

EnvironmentFile=/etc/nexabank/user-service.env
ExecStart=/usr/bin/node src/index.js

Restart=on-failure
RestartSec=10

StandardOutput=journal
StandardError=journal
SyslogIdentifier=nexabank-user-service

PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/etc/nexabank /tmp

[Install]
WantedBy=multi-user.target
